# ansible/config/supervisord/confd.ini
#
#   Author: Rohith
#   Date: 2015-06-11 11:24:32 +0100 (Thu, 11 Jun 2015)
#
#  vim:ts=2:sw=2:et
#
[program:confd]
command=/opt/bin/confd -backend=etcd -node=%(ENV_ETCD_HOSTS)s -watch=true
