# bastion/config/sshd.ini
#
#   Author: Rohith
#   Date: 2015-06-04 14:36:18 +0100 (Thu, 04 Jun 2015)
#
#  vim:ts=2:sw=2:et
#
[program:sshd]
command=/usr/sbin/sshd -p 2222 -D

